# systemd/executor.service
[Unit]
Description=Executor (Build Zero 0d) — paper-mode daemon
Documentation=file:///root/trading-wiki/specs/0d-executor.md
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=root
WorkingDirectory=/root/executor
Environment=PYTHONUNBUFFERED=1
Environment=EXECUTOR_LOG_LEVEL=INFO
Environment=EXECUTOR_AUDIT_DIR=/root/executor/audit-logs/paper_live
Environment=PAPER_MODE=true
# Phase 5 TODO: remove once venue adapter wires real orderbook subscriptions.
# Paper daemon has no orderbook provider — this flag allows the liquidity gate
# to be skipped. Without a real provider the gate silently approves, which is
# acceptable for synthetic-quote paper runs but NOT for real capital.
Environment=EXECUTOR_PAPER_MODE_NO_ORDERBOOK=true
ExecStart=/root/executor/.venv/bin/python -m executor --daemon

# Restart policy — crash-loop guard: give up after 3 restarts in 60s so
# systemd doesn't mask a broken pipeline the way Phase 4's burn-in did.
Restart=on-failure
RestartSec=5
StartLimitIntervalSec=60
StartLimitBurst=3

# Graceful shutdown: SIGTERM triggers bus drain + audit flush + STATE_SAVED.
KillSignal=SIGTERM
TimeoutStopSec=30

# Append stdout/stderr to a single log file; journalctl still sees it via
# systemd's stdlib wrapper.
StandardOutput=append:/var/log/executor.log
StandardError=append:/var/log/executor.log

# Resource caps.
MemoryMax=1G
LimitNOFILE=65536

[Install]
WantedBy=multi-user.target
